# deploy/lego-cam.service
[Unit]
Description=lego-cam sensor-triggered recorder
After=multi-user.target

[Service]
Type=simple
# Use the same user that runs Thonny / your Pi login (e.g. machon). Must be in group i2c for TMF8820.
User=machon
Group=machon

# Project directory (change if your code is elsewhere)
WorkingDirectory=/home/machon/Downloads/legocam/lego-cam-1

# Run from project dir so "src.lego_cam.main" works. Uses project config (same as Thonny).
ExecStart=/usr/bin/python3 -m src.lego_cam.main --config config.example.toml --log-level INFO

Restart=always
RestartSec=5

StandardOutput=journal
StandardError=journal

# Allow camera and I2C
PrivateTmp=false

[Install]
WantedBy=multi-user.target
